# root/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/with-contenv bash

if [ ! "${GOACCESS_ENABLED}" = "true" ]; then
    s6-svc -d .
    exit 0
fi

exec s6-notifyoncheck -d -n 0 \
    /usr/sbin/nginx -g 'daemon off;'
