# root/etc/services.d/traefik/finish
#!/usr/bin/with-contenv bash

rm -rf \
    /tmp/traefik.*
